# bitcoin_trader/app/deployment/bitcoin_trader.service
[Unit]
Description=Bitcoin High Frequency Trading Strategy
After=network.target

[Service]
Type=simple
User=trader
WorkingDirectory=/home/trader/bitcoin_trader
Environment="PATH=/home/trader/bitcoin_trader/.venv/bin:/usr/local/bin:/usr/bin:/bin"

# 启动命令
ExecStart=/home/trader/bitcoin_trader/.venv/bin/python run/high_frequency.py --mode live

# 自动重启配置
Restart=always
RestartSec=10

# 重启限制（10分钟内最多重启5次）
StartLimitInterval=600
StartLimitBurst=5

# 日志配置
StandardOutput=append:/home/trader/bitcoin_trader/logs/service.log
StandardError=append:/home/trader/bitcoin_trader/logs/service_error.log

# 资源限制
LimitNOFILE=65536
MemoryLimit=2G

# 优雅关闭
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
